# packaging/systemd/vdsmd.service
[Unit]
Description=Virtual Desktop Server Manager (Rust rewrite)
Documentation=https://ovirt.org
Requires=libvirtd.service
After=libvirtd.service supervdsmd.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vdsm
Group=vdsm
# DMI serial fields (product_serial / chassis_serial / board_serial) are
# mode 0400 in /sys/class/dmi/id — only root can read them. Copy what we
# need into /run/vdsm/ as root (the `+` prefix grants this ExecStartPre
# elevated privileges even though Service.User=vdsm), then vdsmd reads
# from there as the vdsm user.
ExecStartPre=+/bin/sh -c 'for f in product_serial chassis_serial board_serial; do \
    [ -r /sys/class/dmi/id/$f ] && cp /sys/class/dmi/id/$f /run/vdsm/$f 2>/dev/null; \
done; chmod 0444 /run/vdsm/product_serial /run/vdsm/chassis_serial /run/vdsm/board_serial 2>/dev/null; true'
ExecStart=/usr/bin/vdsmd
Environment=SSLKEYLOGFILE=/var/log/vdsm/sslkeys.log
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
RuntimeDirectory=vdsm
StateDirectory=vdsm
LogsDirectory=vdsm

# Hardening — adjust as we wire real subsystems.
# NoNewPrivileges intentionally OFF for now: vdsmd shells out to
# `sudo mount` (sudoers-gated to /usr/bin/mount + umount only) for
# NFS mounts during connectStorageServer. Real upstream vdsm uses
# supervdsmd over a unix socket for privileged ops; once we ship
# that, this can flip back on.
#NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
# Mount namespace propagation: PrivateTmp=true sets up a private mount
# namespace whose mounts default to slave propagation — meaning NFS
# mounts made by `sudo mount` inside the unit don't reach the host
# namespace, and vanish when the unit restarts. `shared` makes mount
# events propagate both directions, so NFS exports stay mounted across
# vdsmd restarts.
MountFlags=shared
ReadWritePaths=/var/lib/vdsm /var/log/vdsm /run/vdsm /rhev/data-center/mnt

[Install]
WantedBy=multi-user.target
