# packaging/systemd/stubs/openvswitch.service
# vdsm-rs stub — host-deploy expects this service to start; real impl
# not present (vdsm-rs v0 drops OVS/OVN). Always succeeds.
[Unit]
Description=openvswitch (vdsm-rs stub)
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
[Install]
WantedBy=multi-user.target
